# systemd/discord-reminder-bot.service
[Unit]
Description=Discord Auto Message Control Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/zacha/repos/Discord_Auto_Message
EnvironmentFile=/home/zacha/repos/Discord_Auto_Message/systemd/discord-auto-message.env
ExecStart=/home/zacha/repos/Discord_Auto_Message/.venv/bin/python /home/zacha/repos/Discord_Auto_Message/src/bot_control.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
